# .config/crontab
*/3 * * * * ~/.config/scripts/offlineimap.sh
