# .config/crontab
*/3 * * * * ~/.local/bin/offlineimap.sh
